# plists/com.hyperliquid.daemon.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS launchd plist for the trading daemon.
     Replace $AGENT_DIR with your actual install path, e.g.:
       /Users/yourname/hyperliquid-agent
     Then copy to ~/Library/LaunchAgents/ and load with:
       launchctl load ~/Library/LaunchAgents/com.hyperliquid.daemon.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.hyperliquid.daemon</string>

	<key>ProgramArguments</key>
	<array>
		<string>$AGENT_DIR/.venv/bin/python</string>
		<string>-m</string>
		<string>cli.main</string>
		<string>daemon</string>
		<string>start</string>
		<string>--tier</string>
		<string>watch</string>
		<string>--mainnet</string>
		<string>--tick</string>
		<string>120</string>
	</array>

	<key>WorkingDirectory</key>
	<string>$AGENT_DIR</string>

	<key>StandardOutPath</key>
	<string>$AGENT_DIR/data/daemon/daemon_launchd.log</string>

	<key>StandardErrorPath</key>
	<string>$AGENT_DIR/data/daemon/daemon_launchd_err.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>$AGENT_DIR</string>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
