# plists/com.hyperliquid.telegram.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- macOS launchd plist for the Telegram bot service.
     Replace $AGENT_DIR with your actual install path, e.g.:
       /Users/yourname/hyperliquid-agent
     Then copy to ~/Library/LaunchAgents/ and load with:
       launchctl load ~/Library/LaunchAgents/com.hyperliquid.telegram.plist
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.hyperliquid.telegram</string>

	<key>ProgramArguments</key>
	<array>
		<string>$AGENT_DIR/.venv/bin/python</string>
		<string>-m</string>
		<string>telegram.bot</string>
	</array>

	<key>WorkingDirectory</key>
	<string>$AGENT_DIR</string>

	<key>StandardOutPath</key>
	<string>$AGENT_DIR/data/daemon/telegram_bot.log</string>

	<key>StandardErrorPath</key>
	<string>$AGENT_DIR/data/daemon/telegram_bot_err.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PYTHONPATH</key>
		<string>$AGENT_DIR</string>
	</dict>

	<key>RunAtLoad</key>
	<false/>

	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
